# meta-xilinx-core/recipes-graphics/x11-common/xserver-nodm-init/xserver-nodm.service.in
[Unit]
Description=Xserver startup without a display manager
StartLimitBurst=5
StartLimitIntervalSec=100

[Service]
EnvironmentFile=/etc/default/xserver-nodm
User=@USER@
ExecStart=/etc/xserver-nodm/Xserver
Restart=always

[Install]
Alias=display-manager.service
