# meta-xilinx-core/recipes-bsp/dfx-mgr/dfx-mgr/dfx-mgr.service
[Unit]
Description=dfx-mgr

[Service]
ExecStart=/usr/bin/dfx-mgr.sh &

[Install]
WantedBy=multi-user.target
