# meta-xilinx-demos/recipes-utils/start-jupyternb/start-jupyternb/jupyter-setup.service
[Unit]
Description=jupyter setup scripts
Requires=network-online.target
After=network-online.target

[Service]
ExecStart=/sbin/start-jupyter.sh

[Install]
WantedBy=default.target
